# home/.hyperlocal/etc/profile.d/z.sh
../../../../vendor/com/github/rupa/z/z.sh